# scripts/90-container-networking.rules
# this runs the script when there is a network change
SUBSYSTEM=="net",           RUN+="/home/ubuntu/docklet/scripts/container-network-limiter.sh"